# systemd/frost-extract-all.timer
[Unit]
Description=Run hourly netcdf file generation from frost data

[Timer]
OnCalendar=*-*-* *:45:00

[Install]
WantedBy=timers.target